# ibossreg.service
[Unit]
Description=iBoss Registration Agent

[Service]
ExecStart=/usr/sbin/ibossreg
Type=simple

[Install]
WantedBy=multi-user.target
